# salt/osgvo-docker-pilot/osgvo-docker-pilot.cron

*/2 * * * *  root  /usr/sbin/osgvo-check-shutdown

