# salt/docker/docker-maintain.cron


  12 * * * *   root   docker system prune --all --force >/dev/null 2>&1

